# com.divoom.plantmonitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.divoom.plantmonitor</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/southnet-mac-server/PROJECTS/divoom/venv/bin/python</string>
        <string>/Users/southnet-mac-server/PROJECTS/divoom/main.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/southnet-mac-server/PROJECTS/divoom</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>60</integer>

    <key>StandardOutPath</key>
    <string>/Users/southnet-mac-server/PROJECTS/divoom/logs/stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/southnet-mac-server/PROJECTS/divoom/logs/stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
